# crontab
* * * * * root /usr/local/bin/python /sync.py >> /dev/stdout 2>&1
